# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:backend]
command=python -m uvicorn backend.main:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips "*" --timeout-keep-alive 120 --log-level info
directory=/app
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.error.log
autorestart=true
startretries=3
exitcodes=0

[program:mcp]
command=python -m uvicorn backend.mcp.server:app --host 0.0.0.0 --port 8001
directory=/app
stdout_logfile=/var/log/supervisor/mcp.log
stderr_logfile=/var/log/supervisor/mcp.error.log
autorestart=true
startretries=3
exitcodes=0

[program:nginx]
command=nginx -g daemon off;
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.error.log
autorestart=true
